# deploy/douyin-monitor.service
[Unit]
Description=Douyin Barrage Monitor
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/opt/douyin-barrage
ExecStart=/opt/douyin-barrage/venv/bin/python app.py --host 0.0.0.0 --port 8080 --auto-start
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=douyin-barrage

# Security hardening
PrivateTmp=yes
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/opt/douyin-barrage/data /opt/douyin-barrage/logs
ReadOnlyPaths=/opt/douyin-barrage

[Install]
WantedBy=multi-user.target
